# src/main/assembly/dist/bin/easy-example-module.service
[Unit]
Description=EASY Example Module Service

[Service]
ExecStart=/opt/dans.knaw.nl/easy-example-module/bin/easy-example-module-service-start

User=easy-example-module
Group=easy-example-module

[Install]
WantedBy=multi-user.target
